# copilot-key-handler.service
[Unit]
Description=Copilot Key Handler - Opens Claude Code on Copilot key press
After=graphical.target

[Service]
Environment=PATH=/usr/bin:/bin:/usr/sbin:/sbin
ExecStart=/usr/bin/python3 -u /usr/lib/copilot-key-handler/main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
